# debian/postinst
#!/bin/bash
# Post-installation script for sparkmon

set -e

# Ensure sparkmon group and user exist for multi-user setups
if ! getent group sparkmon >/dev/null 2>&1; then
    addgroup --system sparkmon || true
fi

if ! getent passwd sparkmon >/dev/null 2>&1; then
    adduser --system --ingroup sparkmon --home /var/lib/.sparkmon sparkmon || true
fi

# Create and set permissions on data directory  
mkdir -p /var/lib/sparkmon
chown sparkmon:sparkmon /var/lib/sparkmon
chmod 750 /var/lib/sparkmon

# Create systemd service for sparkmon (optional, for background monitoring)
if [ -f /etc/systemd/system/sparkmon.service ]; then
    systemctl daemon-reload || true
    systemctl enable sparkmon || true
fi

# Install shell completion if available
mkdir -p /etc/bash_completion.d/
if [ -f /usr/share/bash-completion/completions/sparkmon ]; then
    cp /usr/share/bash-completion/completions/sparkmon /etc/bash_completion.d/sparkmon || true
fi

echo "Sparkmon installed successfully!"
echo "Configure nodes in /etc/sparkmon/config.yaml or run:"
echo "  sparkmon up --help"
echo "  sparkmon nodes add name=host"
echo ""
echo "For cluster setup, use:"
echo "  sparkmon deploy user@node1,user@node2"